# systemd/user/change_wallpaper.service
[Unit]
Description=Alterar papel de parede usando SWWW
After=swww.service

[Service]
ExecStart=/bin/bash ~/.config/swww/swww.sh

[Install]
WantedBy=default.target